# config/supervisor/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm -F
redirect_stderr=true
